# bo2sale.service
[Unit]
Description=Telegram Bot for Home Marketplace (bo2sale)
After=network.target

[Service]
Type=simple
WorkingDirectory=/root/bo2sale
ExecStart=/usr/bin/python3 /root/bo2sale/bot.py
Restart=on-failure
RestartSec=5
User=root

[Install]
WantedBy=multi-user.target
